# nre/ansible/roles/iota-node/files/iota-node.service
[Unit]
Description=Iota Node

[Service]
User=iota
WorkingDirectory=/opt/iota/
Environment=RUST_BACKTRACE=1
Environment=RUST_LOG=info,iota_core=debug,narwhal=debug,consensus=debug,jsonrpsee=error
ExecStart=/opt/iota/bin/iota-node --config-path /opt/iota/config/validator.yaml
Restart=always

[Install]
WantedBy=multi-user.target
